# runtime/rescue-data-volume.service
[Unit]
Description=Prepare the rescue-web application data volume
After=local-fs.target
Before=rescue-web.service

[Service]
Type=oneshot
ExecStart=/bin/bash /opt/lab/drills/fixtures/02-full-filesystem.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
